# files/etc/supervisor/conf.d/supervisord.conf
[program:gunicorn]
command=/usr/local/bin/gunicorn-start.sh

[program:nginx]
command=/usr/sbin/nginx -g 'daemon on; master_process on;'
